# misc_scripts/cowrie.socket
[Unit]
Description=Cowrie Socket
PartOf=cowrie.service

[Socket]
ListenStream= 22

[Install]
WantedBy=sockets.target